# docker/supervisord.conf
; Quant-AI Dashboard v2.1.1 single-image process layout
; Nginx serves the static frontend and proxies the API to Uvicorn.
; The daemon process can be enabled or disabled through ENABLE_DAEMON.

[supervisord]
nodaemon=true
user=root
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=5MB
logfile_backups=1
loglevel=warn
minfds=1024
minprocs=200

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uvicorn]
command=python -m uvicorn api.main:app --host 0.0.0.0 --port 8685 --workers 1 --limit-concurrency 8 --lifespan on --no-access-log
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=10
priority=20
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:daemon]
command=python run_daemon.py
directory=/app
autostart=%(ENV_ENABLE_DAEMON)s
autorestart=true
startretries=3
startsecs=5
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
